# gfp-pckmgr-updater.service
[Unit]
Description=GFP Package Manager Update Checker
After=network.target
Requires=gfp-pckmgr.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gfp-pckmgr
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/gfp-pckmgr"
ExecStart=/usr/bin/python3 /opt/gfp-pckmgr/check_updates.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=gfp-pckmgr-updater

[Install]
WantedBy=multi-user.target 